# linux/ubuntu/root/etc/profile.d/oraclejdk.sh
#export J2SDKDIR=/usr/lib/jvm/oracle_jdk7
#export J2REDIR=/usr/lib/jvm/oracle_jdk7/jre
#export PATH=$PATH:/usr/lib/jvm/oracle_jdk7/bin:/usr/lib/jvm/oracle_jdk7/db/bin:/usr/lib/jvm/oracle_jdk7/jre/bin
#export JAVA_HOME=/usr/lib/jvm/oracle_jdk7
#export DERBY_HOME=/usr/lib/jvm/oracle_jdk7/db

# Java 8
export IDEA_JDK=/usr/lib/jvm/oracle_jdk8
# SHOULD USE THIS!!! export IDEA_JDK=/<path>/idea-IU-181.4445.78/jre64
#export J2SDKDIR=/usr/lib/jvm/oracle_jdk8
#export J2REDIR=/usr/lib/jvm/oracle_jdk8/jre
export PATH=$PATH:/usr/lib/jvm/oracle_jdk8/bin:/usr/lib/jvm/oracle_jdk8/db/bin:/usr/lib/jvm/oracle_jdk8/jre/bin
export JAVA_HOME=/usr/lib/jvm/oracle_jdk8
#export DERBY_HOME=/usr/lib/jvm/oracle_jdk8/db
